# deploy/vote-bot.service.example
[Unit]
Description=Telegram Vote Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vote-bot
Group=vote-bot
WorkingDirectory=/opt/vote-bot
EnvironmentFile=/opt/vote-bot/.env
ExecStart=/usr/bin/yarn start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
